# scripts/opencode.service
[Unit]
Description=OpenCode Serve
After=network.target

[Service]
Type=simple
User=user1
WorkingDirectory=/home/user1/CollegeLMS
ExecStart=/home/user1/.local/bin/opencode serve --port 4096
Restart=always
RestartSec=5
Environment=HOME=/home/user1

[Install]
WantedBy=multi-user.target
