# scripts/opencode.service
[Unit]
Description=OpenCode Serve
After=network.target

[Service]
Type=simple
User=opencode
WorkingDirectory=/home/opencode/CollegeLMS
ExecStart=/home/opencode/.local/bin/opencode serve --port 4096
Restart=always
RestartSec=5
Environment=HOME=/home/opencode

[Install]
WantedBy=multi-user.target
